# crontab
*/30 * * * * /bin/bash /simulator/sim_autoexecutor.sh >> /simulator/cron.log 2>&1

