# crontab
*/5 * * * * /simulator/sim_autoexecutor.sh
